# packages/gcid/gcid-140.tar.gz/gcid-140/files/gcid.service
[Unit]
Description=Medicine for care laws, Poison for genocide laws
After=multi-user.target

[Service]
#Type=simple
User=gcid
Group=gcid
DynamicUser=yes
StateDirectory=gcid
LogsDirectory=gcid
CacheDirectory=gcid
ExecStart=/usr/local/bin/gcided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
